# crontab
* * * * * python {path-to-project}/cron/disconnect-idle-users.py